# xinetd.d/thp-ftpd
# default: on
# description: thp-ftpd calls the generic thpsvcs with param "ftp",
#	resulting in an ftpd emulation.

service thp-ftp
{
	type			= UNLISTED
	socket_type		= stream
	protocol		= tcp
	port			= 40021
	wait			= no
	user			= nobody
	server			= /usr/local/thp/logthis
	server_args		= ftp
	nice			= 10
	disable			= yes
        instances               = 10
	per_source		= 1

}
